# celery_worker.ini
[program:celery_worker]
command=/home/surajfc/news/newsenv/bin/celery -A news worker -l info
directory=/home/surajfc/news
user=surajfc
numprocs=1
stdout_logfile=/var/log/celery/celery_worker.log
stderr_logfile=/var/log/celery/celery_worker.log
autostart=true
autorestart=true
startsecs = 10
stopwaitsecs = 600
priority=999
